# reference_skills/siem-ticket-bridge/deploy/systemd/siem-ticket-bridge.service
[Unit]
Description=SIEM-to-Ticketing Bridge
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{DEPLOY_DIR}}
EnvironmentFile={{DEPLOY_DIR}}/.env
ExecStart=/usr/bin/python3 -m siem_ticket_bridge.bridge
Restart=on-failure
RestartSec=30
StandardOutput=append:{{LOG_DIR}}/bridge.log
StandardError=append:{{LOG_DIR}}/bridge-error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{LOG_DIR}} {{STATE_DIR}}

[Install]
WantedBy=multi-user.target
